# etc/periodic/weekly/geoip-update
#!/bin/sh
#

wkit geoip-get
wkit nginx reload
